# scripts/claw-cli-scheduler.timer
[Unit]
Description=Claw CLI Scheduler Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
